# syslog-elk-stack/curator/periodic/daily/curator
#!/bin/sh
curator --config /etc/curator/curator.yml /etc/curator/actions.yml
